# build_workspace/BUILD/forge/libinput-rs.service
[Service]
Name = "libinput-rs"
Description = "Memory-safe input daemon for Wayland"
Type = "daemon"
ExecStart = "/usr/local/bin/libinput-rs"
User = "Sisyphus"
Group = "input"

[Dependencies]
Requires = ["udev.service"]
After = ["udev.service"]
Before = ["display-manager.service"]

[Restart]
Condition = "always"
Delay = "2s"
