# systemd/libinput-rs.service
[Unit]
Description=Rust Input Management Daemon
DefaultDependencies=no
After=systemd-udevd.service local-fs.target
Before=display-manager.service
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/libinput-rs
Restart=always
RestartSec=1s
User=root
Group=input
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/dev/input /dev/uinput /etc/libinput-rs
CapabilityBoundingSet=CAP_SYS_ADMIN

[Install]
WantedBy=graphical.target
